# ifi/cron/crontab
MAILTO="henrste@ifi.uio.no,matsasc@ifi.uio.no"

*/5 * * * * ~/drift/ifi/check-dropbox.sh
0 9 * * * ~/drift/ifi/rotate-backups/rotate-backups.sh
